# .config/systemd/user/audio-autoswitch.service
[Unit]
Description=Auto-switch audio output to highest priority device
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/audio-autoswitch.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
